# etc/supervisor/conf.d/exim4.conf
[program:exim4]
command=/usr/sbin/exim -bd -v