# etc/supervisor/conf.d/fcgiwrap.conf
[program:fcgiwrap]
command=fcgiwrap -s unix:/tmp/fcgiwrap.socket
user=www-data